# rootfs/etc/services.d/mqtt-bridge/run
#!/bin/sh
exec /usr/bin/mqtt_to_stdin.sh
