# rootfs/etc/services.d/wmbus_mqtt_bridge/run
#!/usr/bin/with-contenv sh
set -e

echo "RUN.SH STARTED"
exec /usr/bin/run.sh
